# system/showmaster.service
[Unit]
Description=ShowMaster V2 — Live show control for RocketShow
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__INSTALL_DIR__
# Use the exact Node binary detected at install time (nvm, fnm, or system)
Environment=PATH=__NODE_BIN_DIR__:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=SM_PORT=3000
ExecStart=__NODE_BIN_DIR__/node server/src/index.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
